# eink-clock.timer
[Unit]
Description=Run clock oneshot every minute

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
